# overlay.d/05rhcos/usr/lib/systemd/system/coreos-update-ca-trust.service.d/run-every-boot.conf
# For the MCO we want to support adding/removing certs dynamically
[Unit]
ConditionFirstBoot=
