# overlay.d/15rhcos-journald-backcompat/usr/lib/systemd/system/systemd-journald.service.d/rhel8-compat.conf
# See https://issues.redhat.com/browse/LOG-3832
# This allows containers running rhel8 to read our journal
[Service]
Environment=SYSTEMD_JOURNAL_COMPACT=0
Environment=SYSTEMD_JOURNAL_KEYED_HASH=0
Environment=SYSTEMD_JOURNAL_COMPRESS=LZ4
